# systemd/otlp-mcp.service
# OTLP MCP Server - Systemd User Unit
#
# This unit runs otlp-mcp as a persistent HTTP service, auto-discovering
# file sources from the otel-collector configuration.
#
# Installation:
#   ln -s ~/src/otlp-mcp/systemd/otlp-mcp.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now otlp-mcp
#
# Verify:
#   systemctl --user status otlp-mcp
#   journalctl --user -u otlp-mcp -f
#
# MCP Endpoint:
#   http://127.0.0.1:4380/mcp
#
# Claude Code Configuration (~/.claude.json):
#   {
#     "mcpServers": {
#       "otlp-mcp": {
#         "type": "streamableHttp",
#         "url": "http://127.0.0.1:4380/mcp"
#       }
#     }
#   }
#
# Customization:
#   --http-host     Bind address (default: 127.0.0.1)
#   --http-port     Listen port (default: 4380)
#   --otel-config   Path to otel-collector config.yaml
#   --file-source   Manual directory for OTLP JSONL files (can repeat)
#   --verbose       Enable detailed logging
#
# Without --otel-config, the OTLP gRPC receiver starts automatically.
# With --otel-config, only file sources are used (agent can enable OTLP via tools).

[Unit]
Description=OTLP MCP Server (HTTP Transport)
After=network.target otelcol.service

[Service]
Type=exec
# NOTE: Adjust --otel-config path to match your otel-collector configuration.
# Common locations:
#   System: /etc/otelcol/config.yaml (shown below)
#   User:   %h/.config/otelcol/config.yaml
# Remove --otel-config entirely to use OTLP gRPC receiver instead of files.
ExecStart=%h/src/otlp-mcp/otlp-mcp serve \
    --transport http \
    --http-host 127.0.0.1 \
    --http-port 4380 \
    --webui-host 100.83.138.103 \
    --webui-port 4381 \
    --otel-config /etc/otelcol/config.yaml
Restart=on-failure
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=otlp-mcp
NoNewPrivileges=true

[Install]
WantedBy=default.target
